# infra/systemd/x0tta-vpn-boot-validate.service
[Unit]
Description=x0tta6bl4 VPN post-boot validation
Documentation=file:/mnt/projects/docs/runbooks/NL_VPN_HEALTH.md
Requires=x0tta-vpn-route-bypass.service x0tta-vpn-watchdog.service
After=network-online.target x0tta-vpn-route-bypass.service x0tta-vpn-watchdog.service
Wants=network-online.target x0tta-node.service

[Service]
Type=oneshot
WorkingDirectory=/mnt/projects
Environment=PROJECT_DIR=/mnt/projects
Environment=VPN_BOOT_VALIDATE_LOG_DIR=/var/log/x0tta6bl4
Environment=VPN_BOOT_VALIDATE_ATTEMPTS=12
Environment=VPN_BOOT_VALIDATE_SLEEP_SEC=10
ExecStart=/usr/local/sbin/x0tta-vpn-boot-validate
