# src-tauri/debian/prerm
#!/bin/sh
set -e

SERVICES="x0tta6bl4-full-core-api.service x0tta6bl4-core-api.service x0tta6bl4-node.service"

case "$1" in
  remove|deconfigure)
    if command -v systemctl >/dev/null 2>&1 && [ -d /run/systemd/system ]; then
      systemctl stop $SERVICES >/dev/null 2>&1 || true
      systemctl disable $SERVICES >/dev/null 2>&1 || true
    fi
    ;;
esac

exit 0
